# scripts/content-generator/launchd/com.localizeshots.content-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localizeshots.content-morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ivorobyev/projects/localizeshots/scripts/content-generator/run-morning.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ivorobyev/projects/localizeshots/scripts/content-generator/logs/launchd-morning.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ivorobyev/projects/localizeshots/scripts/content-generator/logs/launchd-morning-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ivorobyev/projects/localizeshots/scripts/content-generator</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
